# meta-edison-distro/recipes-support/crashlog/files/crashlog.service
[Unit]
Description=Crashlog service

[Service]
ExecStart=/usr/sbin/retrieve_crashlog.sh
Restart=no

[Install]
WantedBy=basic.target

