# meta-intel-edison-distro/recipes-support/cleanjournal/files/cleanjournal.service
[Unit]
Description=Cleanjournal service

[Service]
ExecStart=/usr/sbin/clean_journal.sh
Restart=no

[Install]
WantedBy=basic.target

