# meta-intel-edison-distro/recipes-support/blink-led/files/blink-led.service
[Unit]
Description=Edison Arduino board LED Blinker

[Service]
ExecStart=/usr/bin/blink-led
ExecStop=/usr/bin/blink-led
RemainAfterExit=Yes

[Install]
WantedBy=basic.target
