# meta-intel-edison-distro/recipes-support/blink-led/files/blink-led.service
[Unit]
Description=Edison Arduino board LED Blinker

[Service]
ExecStart=/usr/bin/blink-led
Restart=on-failure

[Install]
WantedBy=basic.target

