# meta-intel-edison-distro/recipes-support/blink-led/files/init-led.service
[Unit]
Description=Edison Arduino board LED Blinker
Before=post-install.service

[Service]
ExecStart=/usr/bin/init-led

[Install]
WantedBy=basic.target
