# meta-intel-edison-distro/recipes-protocols/xl2tpd/files/xl2tpd.service
[Unit]
Description=Layer 2 Tunneling Protocol Daemon
After=network.target

[Service]
ExecStart=/usr/sbin/xl2tpd -c /etc/xl2tpd.conf  -D
Restart=on-failure

[Install]
WantedBy=multi-user.target

