# meta-intel-edison-bsp/recipes-support/sketch-check/files/sketch-check.service
[Unit]
Description=Edison sketch check service

[Service]
ExecStart=/usr/bin/sketch-check.sh 
Restart=on-failure

[Install]
WantedBy=default.target

