# meta-intel-edison-distro/recipes-support/battery-voltage/battery-voltage/battery-voltage.service
[Unit]
Description=Battery Voltage daemon

[Service]
ExecStart=/usr/bin/battery-voltage
Restart=on-failure

[Install]
WantedBy=basic.target

